# supervisord.conf
[supervisord]
nodaemon=true
logfile=/workspace/logs/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:vllm-server]
command=bash /services/start_vllm.sh
directory=/workspace
autostart=true
autorestart=true
startsecs=30
startretries=3
priority=10
stdout_logfile=/workspace/logs/vllm.log
stderr_logfile=/workspace/logs/vllm.err
environment=HF_HOME="/workspace/models/hf_cache"

[program:stt-server]
command=python /services/stt_server.py
directory=/workspace
autostart=true
autorestart=true
startsecs=10
startretries=3
priority=15
stdout_logfile=/workspace/logs/stt.log
stderr_logfile=/workspace/logs/stt.err
environment=HF_HOME="/workspace/models/hf_cache",CUDA_VISIBLE_DEVICES="0"

[program:tts-server]
command=python /services/tts_server.py
directory=/workspace
autostart=true
autorestart=true
startsecs=30
startretries=3
priority=15
stdout_logfile=/workspace/logs/tts.log
stderr_logfile=/workspace/logs/tts.err
environment=HF_HOME="/workspace/models/hf_cache",CUDA_VISIBLE_DEVICES="0",COSYVOICE_PATH="/opt/cosyvoice"

[program:api-server]
command=uvicorn app.api.main:app --host 0.0.0.0 --port 8888 --workers 1
directory=/
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=20
stdout_logfile=/workspace/logs/api.log
stderr_logfile=/workspace/logs/api.err

[program:jupyterlab]
command=jupyter lab --ip=0.0.0.0 --port=8889 --no-browser --allow-root --NotebookApp.token='' --NotebookApp.password=''
directory=/workspace
autostart=true
autorestart=true
priority=30
stdout_logfile=/workspace/logs/jupyter.log
stderr_logfile=/workspace/logs/jupyter.err
